# scripts/com.maxspevack.fishwrap.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maxspevack.fishwrap.refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/max/gemini/fishwrap/scripts/refresh_demos.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>21600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/max/gemini/fishwrap/logs/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/max/gemini/fishwrap/logs/launchd.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/max/gemini/fishwrap</string>
</dict>
</plist>
